# scallop-reset.service
[Unit]
Description=Reset options for Scallop systems

[Service]
ExecStart=/usr/bin/scallop --reset
